# configs/services/initCheckConfigurationUpdated.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/scripts/check_configuration_updated.sh

[Install]
WantedBy=default.target